# etc/init/gitlab-webhook-receiver.conf
start on runlevel [2345]
stop on shutdown
respawn
exec su -s /bin/sh -c 'HOME=/var/lib/gitlab-webhook-receiver exec "$0" "$@" &>/dev/null' gitlab-webhook-receiver -- gitlab-webhook-receiver 
